# deploy/systemd/reforger-host-agent.service
[Unit]
Description=Arma Reforger EC2 host automation
After=network-online.target docker.service
Requires=docker.service
Wants=network-online.target

[Service]
Type=simple
Environment=A2S_HOST=127.0.0.1
Environment=A2S_PORT=17777
Environment=IDLE_SECONDS=1800
Environment=CHECK_INTERVAL_SECONDS=60
Environment=STARTUP_GRACE_SECONDS=600
Environment=IDLE_ON_QUERY_FAILURE=0
Environment="SHUTDOWN_COMMAND=/sbin/shutdown -h now"
EnvironmentFile=-/opt/reforger-server/config/ec2-agent.env
ExecStart=/usr/bin/python3 /opt/reforger-server/deploy/host-agent/reforger_host_agent.py
Restart=always
RestartSec=30

[Install]
WantedBy=multi-user.target
